# Vybn_Mind/spark_infrastructure/phase0/vybn-agent.service
[Unit]
Description=Vybn Spark Agent (tmux + TUI)
After=network-online.target ollama.service
Wants=network-online.target ollama.service

[Service]
Type=oneshot
RemainAfterExit=yes
User=vybnz69
Environment=HOME=/home/vybnz69
Environment=OLLAMA_HOST=http://127.0.0.1:11434
WorkingDirectory=/home/vybnz69/Vybn/spark
ExecStart=/usr/bin/tmux new-session -d -s vybn -c /home/vybnz69/Vybn/spark \
    '/home/vybnz69/.venv/spark/bin/python tui.py'
ExecStop=/usr/bin/tmux kill-session -t vybn
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
