# spark/systemd/vybn-vllm.service
# User-level wrapper around launch-cluster.sh. Keeps the vLLM Ray cluster
# (Nemotron 120B across both Sparks) up. Ray distributed backend is required
# because world_size=2 with one GPU per node.
# Source of truth: ~/Vybn/spark/systemd/.
[Unit]
Description=Vybn vLLM cluster (Nemotron 120B, 2 Sparks via Ray)
After=network-online.target
Wants=network-online.target
# Docker is a system-level unit and invisible from user systemd; the daemon
# is always running, so we depend on the binary instead via ConditionPath.
ConditionPathExists=/var/run/docker.sock
StartLimitIntervalSec=600
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=%h/spark-vllm-docker
# Always clean-stop before starting so we never layer on a zombie cluster.
ExecStartPre=/bin/bash -c '%h/spark-vllm-docker/launch-cluster.sh -n 169.254.246.181,169.254.51.101 stop || true'
# Foreground (no -d) so systemd supervises.
ExecStart=/home/vybnz69/spark-vllm-docker/launch-cluster.sh \
  -n 169.254.246.181,169.254.51.101 \
  exec \
  vllm serve nvidia/NVIDIA-Nemotron-3-Super-120B-A12B-FP8 \
  --port 8000 --host 0.0.0.0 \
  --gpu-memory-utilization 0.85 --max-model-len 8192 \
  --tensor-parallel-size 1 --pipeline-parallel-size 2 \
  --distributed-executor-backend ray \
  --trust-remote-code
ExecStop=/bin/bash -c '%h/spark-vllm-docker/launch-cluster.sh -n 169.254.246.181,169.254.51.101 stop || true'
Restart=always
RestartSec=30
# Cold load across 2 Sparks is ~10-13 min (26 FP8 shards @ ~25s + cudagraph).
# Anything shorter than 15 min risks SIGTERM-ing a healthy warm-up.
TimeoutStartSec=1200
TimeoutStopSec=60
StandardOutput=append:%h/logs/vllm.log
StandardError=append:%h/logs/vllm.log

[Install]
WantedBy=default.target
