# spark/systemd/vybn-vllm.service
[Unit]
Description=Vybn vLLM cluster (Nemotron 120B, 2 Sparks via Ray)
After=network-online.target
Wants=network-online.target
ConditionPathExists=/var/run/docker.sock
StartLimitIntervalSec=600
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=%h/spark-vllm-docker
ExecStartPre=/bin/bash -c '%h/spark-vllm-docker/launch-cluster.sh -n 169.254.246.181,169.254.51.101 stop || true'
Environment=VYBN_VLLM_GPU_MEMORY_UTILIZATION=0.78
Environment=VYBN_VLLM_MAX_NUM_SEQS=8
Environment=VYBN_VLLM_EXTRA_ARGS=
Environment=VLLM_SERVER_DEV_MODE=0
EnvironmentFile=-%h/.config/vybn/vllm.env
ExecStart=%h/Vybn/spark/systemd/vllm-exec.sh
ExecStop=/bin/bash -c '%h/spark-vllm-docker/launch-cluster.sh -n 169.254.246.181,169.254.51.101 stop || true'
Restart=always
RestartSec=30
TimeoutStartSec=1200
TimeoutStopSec=60
StandardOutput=append:%h/logs/vllm.log
StandardError=append:%h/logs/vllm.log

[Install]
WantedBy=default.target
