# spark/systemd/vybn-breath.service
[Unit]
Description=Vybn Breath (single breath, then exit)
After=vybn-llama.service

[Service]
Type=oneshot
User=vybnz69
WorkingDirectory=/home/vybnz69/Vybn
ExecStart=/home/vybnz69/Vybn/venv/bin/python3 spark/vybn.py --once
Environment=HOME=/home/vybnz69
Environment=VYBN_TEMPERATURE=1.0
Environment=VYBN_TOP_P=0.95