# spark/systemd/vybn-walk-daemon.service
# User-level. Source of truth: ~/Vybn/spark/systemd/.
# walk_daemon hard-codes 127.0.0.1:8101 (walk_daemon.py).
[Unit]
Description=Vybn walk daemon (walk_daemon.py on :8101)
After=network-online.target vybn-deep-memory.service
Wants=network-online.target vybn-deep-memory.service
StartLimitIntervalSec=300
StartLimitBurst=20

[Service]
Type=simple
WorkingDirectory=%h/vybn-phase
ExecStartPre=/bin/sh -c 'fuser -k 8101/tcp >/dev/null 2>&1 || true; sleep 1'
ExecStart=/usr/bin/python3 %h/vybn-phase/walk_daemon.py
Restart=always
RestartSec=5
TimeoutStartSec=60
TimeoutStopSec=15
StandardOutput=append:%h/logs/walk_daemon.log
StandardError=append:%h/logs/walk_daemon.log

[Install]
WantedBy=default.target
