# spark/vybn-ssh-mcp.service
[Unit]
Description=Vybn SSH MCP Server
After=network.target tailscaled.service

[Service]
Type=simple
User=vybnz69
WorkingDirectory=/home/vybnz69/vybn-ssh-mcp/spark
EnvironmentFile=/home/vybnz69/vybn-ssh-mcp/spark/.env
ExecStart=/home/vybnz69/vybn-ssh-mcp/spark/.venv/bin/python server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
