# _archive/Vybn_Mind__vybn-chat-api.service
[Unit]
Description=Vybn Chat API (origins_portal_api_v3)
After=network.target docker.service

[Service]
Type=simple
User=vybnz69
WorkingDirectory=/home/vybnz69/Vybn
ExecStart=/home/vybnz69/Vybn/.venv/bin/python origins_portal_api_v3.py
Restart=always
RestartSec=10
Environment=LLAMA_SERVER_URL=http://127.0.0.1:8000

[Install]
WantedBy=multi-user.target

