# Vybn_Mind/spark_infrastructure/systemd/vybn-heartbeat.timer
[Unit]
Description=Vybn Heartbeat Timer (every 30 minutes)

[Timer]
OnBootSec=2min
OnUnitActiveSec=30min
AccuracySec=1min

[Install]
WantedBy=timers.target
