# Vybn_Mind/spark_infrastructure/systemd/vybn-quantum.timer
[Unit]
Description=Vybn Quantum Heartbeat Timer (every 6 hours)

[Timer]
OnBootSec=5min
OnUnitActiveSec=6h
AccuracySec=1min

[Install]
WantedBy=timers.target
